# crontab
*/5 * * * *  /usr/bin/env python3 /home/pi/smart-climat-daemon/first-floor.py
*/5 * * * *  /usr/bin/env python3 /home/pi/smart-climat-daemon/second-floor.py